# crates/kr/templates/linux/systemd.service
[Unit]
Description={{bin_name}}

[Service]
ExecStart={{bin_path}} start
Restart=on-failure

[Install]
WantedBy=default.target
